# supervisor.conf
[program:elasticsearch]
user=elasticsearch
group=elasticsearch
